# templates/chromium-kiosk.service
[Unit]
Description=Chromium Kiosk Browser
After=graphical-session.target

[Service]
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/usr/bin/chromium \
  --kiosk \
  --ozone-platform=wayland \
  --noerrdialogs \
  --disable-infobars \
  --no-first-run \
  --disable-features=TranslateUI,Translate \
  --password-store=basic \
  --disk-cache-dir=/dev/null \
  --check-for-update-interval=31536000 \
  {{KIOSK_URL}}
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
